# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:proxy]
command=python3 /app/proxy/main.py
directory=/data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:api-server]
command=node --enable-source-maps /app/api-server/dist/index.mjs
directory=/app/api-server
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
; Force internal port 8080 — Zeabur injects PORT=<random> for nginx, api-server must stay on 8080
; ADMIN_PASSWORD and SESSION_SECRET are passed through from the container environment
environment=PORT=8080,NODE_ENV=production,DATA_DIR=/data,LOG_FILE=/data/proxy.log,PROXY_INTERNAL_URL=http://127.0.0.1:8000,ADMIN_PASSWORD=%(ENV_ADMIN_PASSWORD)s,SESSION_SECRET=%(ENV_SESSION_SECRET)s

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
